# DEBIAN/postinst
#!/bin/bash
set -e

# Install dependencies
echo "Installing dependencies..."
cd /opt/ssh-globe-tracker/server
npm install --production

# Create systemd service
systemctl daemon-reload
systemctl enable ssh-globe-tracker.service

echo "SSH Globe Tracker installed successfully!"
echo "Use 'ssh-globe-tracker start' to start the service"
echo "Use 'ssh-globe-tracker status' to check status"
echo "Use 'ssh-globe-tracker logs' to view logs"
echo "Access the web interface at: http://localhost:3001"
